# files/supervisord.conf
[supervisord]
nodaemon=true

[program:go]
command=/bin/bash -c "cd /go && go run *.go"
stdout_events_enabled=true
stderr_events_enabled=true
